# deploy/alicloud/ubuntu-22.04/templates/my-first-expo-app-email-agent.service
[Unit]
Description=My First Expo App Email Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=__APP_ROOT__/email-agent-backend
ExecStart=/usr/bin/node __APP_ROOT__/email-agent-backend/dist/app.js
Restart=always
RestartSec=5
Environment=PORT=__EMAIL_AGENT_PORT__
EnvironmentFile=-__APP_ROOT__/email-agent-backend/.env
User=__APP_USER__
Group=__APP_GROUP__

[Install]
WantedBy=multi-user.target
